# root/etc/services.d/cron/run
#!/bin/execlineb -P

foreground { 
  s6-svwait -u /var/run/s6/services/syslog
}

exec -c
fdmove -c 2 1

/usr/sbin/crond -f -S -l 0 -c /etc/crontabs
