# root/etc/services.d/cron/run
#!/bin/execlineb -P

#exec -c
#fdmove -c 2 1

#/usr/sbin/cron
exec /usr/sbin/cron
