# root/etc/cron.d/cheky-check
*/30 * * * * php /var/www/html/check.php
